# C/dist/etc/mark/mark.service
[Unit]
Description=Mark machine manager

[Service]
Type=simple
ExecStart=/etc/mark/mark --start --demo

